# php72-apache-broadcast4/files/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/bin/bash -c '/etc/init.d/ssh start'

[program:php-fpm]
command=/bin/bash -c 'php-fpm --nodaemonize'

[program:apache2]
command=/bin/bash -c "/usr/sbin/apache2ctl -DFOREGROUND"

[program:logstash-forwarder]
command=/bin/bash -c "service logstash-forwarder start"

[include]
files = /etc/supervisor/conf.d/additional_*.conf
